# pi-service/netboot-imager.service
[Unit]
Description=Netboot Image Writer Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/netboot-imager
ExecStart=/usr/bin/python3 /opt/netboot-imager/netboot-imager-service.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
